# dyneII/startup/udev/70-firmware.rules
ACTION=="add", SUBSYSTEM=="firmware", RUN+="/sbin/udev_run_hotplugd"
